# helpers/systemd/monitoring@exporters.service.d/order.conf
[Unit]
Requires=monitoring@proxy.service monitoring@zabbix.service monitoring@promgraf.service
After=monitoring@proxy.service monitoring@zabbix.service monitoring@promgraf.service
