# helpers/systemd/monitoring@zabbix.service.d/order.conf
[Unit]
Requires=monitoring@proxy.service
After=monitoring@proxy.service
