# helpers/systemd/monitoring@promgraf.service.d/order.conf
[Unit]
Requires=monitoring@proxy.service monitoring@zabbix.service
After=monitoring@proxy.service monitoring@zabbix.service
